# lightning-integration.service
[Service]
Type=simple
ExecStart=/opt/lightning-integration-systemd/run.sh
